# ubuntu/root/etc/s6-overlay/s6-rc.d/php-fpm-prepare/run
#!/command/exec bash
cp /etc/s6-overlay/files/www.conf /etc/php/$(ls /etc/php)/fpm/pool.d/www.conf
