# root/etc/s6-overlay/s6-rc.d/welcome-message/run
#!/command/with-contenv bash

. /etc/s6-overlay/lib/logger

print_welcome_page
