# root/etc/s6-overlay/s6-rc.d/php-fpm-prepare/run
#!/command/exec bash

mkdir /run/php-fpm
chown nginx: /run/php-fpm
sed -i -E 's/^(user|group).*/\1 = nagios/' /etc/php-fpm.d/www.conf
