# fedora/root/etc/s6-overlay/s6-rc.d/fcgiwrap-fix/run
#!/command/exec bash
SOCKET=/run/fcgiwrap.sock
chown nginx: ${SOCKET}
chmod 700 ${SOCKET}
