# fedora/root/etc/s6-overlay/s6-rc.d/fcgiwrap/run
#!/command/execlineb -P
/usr/sbin/fcgiwrap -f -c 1 -s unix:/run/fcgiwrap.sock
