# root/etc/s6-overlay/s6-rc.d/fcgiwrap-fix/run
#!/command/exec bash

. /etc/s6-overlay/lib/logger

SOCKET=/run/fcgiwrap.sock

info "Assign fcgiwrap socket to user nginx."

chown nginx: ${SOCKET}
chmod 700 ${SOCKET}
