# ubuntu/root/etc/s6-overlay/s6-rc.d/fcgiwrap/run
#!/command/exec bash
SOCKET=/run/fcgiwrap.socket
/usr/sbin/fcgiwrap -f -c 1 -s unix:${SOCKET}
chown nginx: ${SOCKET}
chmod 700 ${SOCKET}
