# ubuntu/root/etc/s6-overlay/s6-rc.d/fcgiwrap/up
spawn-fcgi -u nginx -g nginx -s /run/fcgiwrap.socket -M 0700 -F 1 -P /run/fcgi.pid -- /usr/sbin/fcgiwrap
